# etc/pve-lxc-syscalld.service.in
[Unit]
Description=Proxmox VE LXC Syscall Daemon
Before=pve-guests.service

[Service]
Type=notify
ExecStart=%LIBEXECDIR%/pve-lxc-syscalld/pve-lxc-syscalld --system /run/pve-lxc-syscalld/socket
RuntimeDirectory=pve-lxc-syscalld
Restart=on-failure

[Install]
WantedBy=multi-user.target
